# airootfs/etc/systemd/system/getty@tty3.service.d/autologin.conf
[Service]
Environment=LANG=en_US.UTF-8
ExecStart=
ExecStart=-/usr/bin/agetty --autologin root --noclear %I $TERM
